# ansible/roles/backup_mongodb/files/backup_db.timer
[Unit]
Description=Backup Alpi DB timer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target